# lib/systemd/system/mdnsbridge.service
[Unit]
Description=mdnsbridge

[Service]
ExecStart=/usr/bin/nmos-mdnsbridge
