# debian/python-mdnsbridge.upstart
description	"nmos-mdnsbridge"
author		"Peter Brightwell <peter.brightwell@bbc.co.uk>"

start on runlevel [2345]
stop on runlevel [016]
respawn

script
	exec sudo -u ipstudio /usr/bin/nmos-mdnsbridge
end script
